# driver-redpanda/deploy/templates/benchmark-code-server.service
[Unit]
Description=Benchmark artifact server
After=network.target

[Service]
WorkingDirectory=/opt/
ExecStart=python3 -m http.server 6000

[Install]
WantedBy=multi-user.target